# svc/appstack-platoon.service
[Unit]
Description=Appstack Platoon service
After=rsyslog.service

[Service]
ExecStart=/siq/svc/platoon/platoon-ctl start
ExecStop=/siq/svc/platoon/platoon-ctl stop
RuntimeDirectory=platoon
Restart=always
KillSignal=SIGQUIT
Type=forking
PIDFile=/var/siq/run/platoon.pid
StandardError=syslog
NotifyAccess=all

[Install]
WantedBy=appstack.service
